# launchd/com.example.feishu-companion-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.feishu-companion-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-dimsu</string>
    <string>/ABSOLUTE/PATH/TO/feishu-companion-bot/.venv/bin/python</string>
    <string>-m</string>
    <string>feishu_companion</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/feishu-companion-bot</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/feishu-companion-bot/bot.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/feishu-companion-bot/bot.err.log</string>
</dict>
</plist>
